# resources/crontab
#  m h dom mon dow user  command
*/15 *   *   *   * root  /srv/archive-auth-mirror/bin/mirror-archive
